# ai-ddr-generator.service
[Unit]
Description=AI DDR Generator - Property Inspection Report Generator
Documentation=file:///opt/ai-ddr-generator/DEPLOYMENT.md
After=network.target

[Service]
Type=oneshot
User=ai-ddr-gen
Group=ai-ddr-gen
WorkingDirectory=/opt/ai-ddr-generator

# Environment
EnvironmentFile=/opt/ai-ddr-generator/.env
Environment="PATH=/opt/ai-ddr-generator/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Execution
ExecStart=/opt/ai-ddr-generator/venv/bin/python /opt/ai-ddr-generator/prod_main.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-ddr-gen

# Security
PrivateTmp=yes
NoNewPrivileges=true
ReadOnlyPaths=/etc
ReadWritePaths=/opt/ai-ddr-generator/outputs
ReadWritePaths=/opt/ai-ddr-generator/data

# Restart policy
Restart=on-failure
RestartSec=60
MaxRetries=3

[Install]
WantedBy=multi-user.target
