# ai-ddr-generator.timer
[Unit]
Description=AI DDR Generator - Daily Run Timer
Documentation=file:///opt/ai-ddr-generator/DEPLOYMENT.md

[Timer]
# Run at 2:00 AM every day
OnCalendar=*-*-* 02:00:00

# Or run every 6 hours
# OnBootSec=1h
# OnUnitActiveSec=6h

# Run immediately if system was down at scheduled time
Persistent=true

[Install]
WantedBy=timers.target
